# supervisord.conf
[supervisord]

nodaemon=true

[program:angularjs]

command=yo angular angularApp && grunt && grunt serve
